# MRL_AI_SYSTEM/particles/boot/systemd/mrl-smartbody.service
[Unit]
Description=MRL particle-smartbody-v2 (body: sense/act/reflex)
Requires=mrl-doctor.service
After=mrl-doctor.service
PartOf=mrl.target

[Service]
Type=simple
User=mrl
Group=mrl
Environment=MRL_BODY_HOST=127.0.0.1
Environment=MRL_BODY_PORT=8787
Environment=MRL_BODY_HUB=http://127.0.0.1:9000
Environment=MRL_BODY_DOCTOR=http://127.0.0.1:8788
Environment=MRL_ACT_WHITELIST=/etc/mrl/act_whitelist.txt
Environment=MRL_TRACE_JOURNAL=/var/lib/mrl/trace/journal.jsonl
ExecStart=/usr/bin/env python3 /opt/mrl/ai-system/MRL_AI_SYSTEM/particles/smartbody_v2/mrl_smartbody.py
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/lib/mrl /var/log/mrl

[Install]
WantedBy=mrl.target
