# MRL_AI_SYSTEM/particles/boot/systemd/mrl-sig-verify.service
[Unit]
Description=MRL particle-sig-verify (L0 signature law)
After=network-online.target
Wants=network-online.target
PartOf=mrl.target

[Service]
Type=simple
User=mrl
Group=mrl
Environment=MRL_SIG_HOST=127.0.0.1
Environment=MRL_SIG_PORT=8801
Environment=MRL_SIG_KEYDIR=/etc/mrl/keys
Environment=MRL_TRACE_JOURNAL=/var/lib/mrl/trace/journal.jsonl
ExecStart=/usr/bin/env python3 /opt/mrl/ai-system/MRL_AI_SYSTEM/particles/sig_verify/mrl_sig_verify.py
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/etc/mrl /var/lib/mrl

[Install]
WantedBy=mrl.target
